# com.zaino.telegram-bot-updater.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zaino.telegram-bot-updater</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>auto-update.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/zaino/zaino-telegram-bot</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/zaino/zaino-telegram-bot/updater.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/zaino/zaino-telegram-bot/updater-error.log</string>
</dict>
</plist>
